# resources/gudgeon.socket
[Unit]
Description=Gudgeon Socket
PartOf=gudgeon.service
After=network.service

[Socket]
ListenStream=0.0.0.0:53
ListenDatagram=0.0.0.0:53
