# resources/gudgeon.socket
[Unit]
Description=Gudgeon Socket

[Socket]
ListenStream=53
NoDelay=true